# etc/inittab
# /etc/inittab: init(8) configuration.
::sysinit:/etc/init.d/rcS

# shell
::askfirst:-/bin/login
